# Server/scripts/rogueserver.service
[Unit]
Description=Rogue Daemon Game Server

[Service]
ExecStart=/home/leif/rogue/go.sh

[Install]
WantedBy=multi-user.target

